# raspberry-pi-files/systemd/smartkart-rfid.service
[Unit]
Description=SmartKart RFID Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/smartkart/raspberry-pi-files
ExecStart=/usr/bin/python3 /home/pi/smartkart/raspberry-pi-files/rfid_service.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables (optional - hardcoded in script)
# Environment="BACKEND_URL=http://172.16.71.172:8001"
# Environment="CART_ID=1234"

[Install]
WantedBy=multi-user.target
